# package/snapcast/snapserver.service
[Unit]
Description=Snapcast server
Wants=avahi-daemon.service
After=network-online.target time-sync.target avahi-daemon.service

[Service]
PermissionsStartOnly=true
RuntimeDirectory=snapserver
RuntimeDirectoryMode=0775
EnvironmentFile=/etc/default/snapserver
ExecStart=/usr/bin/snapserver --logging.sink=system $SNAPSERVER_OPTS
User=snappi
Group=audio
Restart=always
LimitRTPRIO=infinity
LimitMEMLOCK=infinity
LimitRTTIME=infinity

[Install]
WantedBy=multi-user.target

