# package/snappi-server-config/uac2.service
[Unit]
Description=Snappi UAC2 Service

[Service]
Type=oneshot
ExecStart=/usr/bin/uac2.sh

[Install]
WantedBy=multi-user.target


